# roles/jenkins_master/files/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/etc/init.d/ssh start

[program:jenkins]
command=/bin/tini -- /usr/local/bin/jenkins.sh